# deploy/com.trumpstockalerts.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trumpstockalerts.bot</string>

    <!-- Run the always-on launcher (continuous polling loop). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/shonsmac/Apps/trump-stock-alerts/deploy/run.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/shonsmac/Apps/trump-stock-alerts</string>

    <!-- Start at login and keep it alive: relaunch if it ever exits/crashes. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Throttle restarts so a crash-loop can't hammer the machine. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/shonsmac/Apps/trump-stock-alerts/data/bot.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shonsmac/Apps/trump-stock-alerts/data/bot.err.log</string>
</dict>
</plist>
